# src/conf/crontab
#program console

## 监控
*/1 * * * * sleep 10;/usr/local/bin/php /www/lucky/src/public/console/index.php request_uri=/monitor > /dev/null 2>&1
